# docker/crontab
# Crontab for Tiki Crawler - Tự động crawl từ database
# Format: minute hour day month weekday command

# Crawl categories từ database mỗi 3 phút (auto-select theo priority)
*/3 * * * * /usr/local/bin/python /app/src/crawl_from_categories.py >> /app/logs/cron-crawler.log 2>&1

# Heartbeat log every 5 minutes
*/5 * * * * echo "[$(date)] Cron is alive" >> /app/logs/cron-heartbeat.log 2>&1

# Keep a newline at the end
